# com.geminiclonenextweb.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.geminiclonenextweb.server</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/aydinkaya/开发/chatgpt-web-midjourney-proxy/start.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/aydinkaya/开发/chatgpt-web-midjourney-proxy</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/aydinkaya/开发/chatgpt-web-midjourney-proxy/logs/gemini-clone.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/aydinkaya/开发/chatgpt-web-midjourney-proxy/logs/gemini-clone-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>development</string>
        <key>PORT</key>
        <string>1002</string>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>